# live-build/config/includes.chroot/etc/systemd/system/spawn-tor-gateway.service
[Unit]
Description=Spawn Tor Traffic Gateway
After=tor.service
Requires=tor.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/spawn-tor-gateway start
ExecStop=/usr/local/bin/spawn-tor-gateway stop

[Install]
WantedBy=multi-user.target
